# teletask_bridge/rootfs/etc/services.d/teletask/run
#!/usr/bin/with-contenv bashio
# Fail fast
set -euo pipefail

bashio::log.info "Starting Teletask Bridge..."

# If your app expects a config file in /app/config.json, ensure it's present or mount it via options.
# Start the bridge in the foreground so s6 can supervise it:
exec python3 /app/main.py