# teletask_bridge/rootfs/etc/services.d/teletask/run
#!/usr/bin/with-contenv bashio
# Fail fast
set -euo pipefail

# Install dependencies (only first run; safe to re-run)
pip3 install --no-cache-dir gmqtt || true

bashio::log.info "Starting Teletask Bridge..."

CONFIG_PATH="$(bashio::config 'config_path')"

if ! bashio::fs.file_exists "$CONFIG_PATH"; then
    bashio::log.fatal "Config file not found at: $CONFIG_PATH"
    exit 1
fi

bashio::log.info "Using configuration: $CONFIG_PATH"

exec python3 /app/main.py "$CONFIG_PATH"
