# teletask_bridge/rootfs/etc/services.d/teletask/run
#!/usr/bin/with-contenv bashio
# Fail fast
set -euo pipefail

bashio::log.info "Starting Teletask Bridge..."

CONFIG_PATH="$(bashio::config 'config_path')"

if ! bashio::fs.file_exists "$CONFIG_PATH"; then
    bashio::log.fatal "Config file not found at: $CONFIG_PATH"
    exit 1
fi

bashio::log.info "Using configuration: $CONFIG_PATH"

# Use the venv Python we created at build time
exec /opt/venv/bin/python /app/main.py "$CONFIG_PATH"
